# docker/supervisor/queue.conf
[program:queue]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work --tries=3
directory=/application
autostart=true
autorestart=true
numprocs=1
user=root